# supervisord.conf
[unix_http_server]
file=/wls/applications/sql_predict/apps/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9011        ; (ip_address:port specifier, *:port for all iface)
username=root              ; (default is no username (open server))
password=Pa888888          ; (default is no password (open server))

[supervisord]
logfile=/wls/applications/sql_predict/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=150MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug                ; (log level;default info; others: debug,warn,trace)
pidfile=/wls/applications/sql_predict/apps/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///wls/applications/sql_predict/apps/supervisor.sock ; use a unix:// URL  for a unix socket

[program:pallas-web]
command=python2.7 -m SimpleHTTPServer 8012
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/wls/applications/sql_predict/apps/web/dist             ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
stopsignal=QUIT               ; signal used to kill process (default TERM)
user=dbcm                   ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/wls/applications/sql_predict/logs/pallas-web.log        ; stdout log path, NONE for none; default AUTO
#environment=APP_ENV="ProductionAdmin"       ; process environment additions (def no adds)
stopasgroup=true
killasgroup=true

[program:pallas-app]
command=/DATA1/virtualenv/ai_python3_virtualenv/bin/python3.6 /wls/virtualenv/ai_python3_virtualenv/bin/gunicorn sqlreview.wsgi:application -t 1800 -w 2 --bind 0.0.0.0:8011
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/wls/applications/sql_predict/apps             ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
stopsignal=QUIT               ; signal used to kill process (default TERM)
user=dbcm                   ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/wls/applications/sql_predict/logs/pallas-app.log        ; stdout log path, NONE for none; default AUTO
#environment=APP_ENV="ProductionAdmin"       ; process environment additions (def no adds)
stopasgroup=true
killasgroup=true

[program:pallas-celery]
command=/DATA1/virtualenv/ai_python3_virtualenv/bin/python3.6 /wls/virtualenv/ai_python3_virtualenv/bin/celery -A sqlreview worker -l info -Ofair
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/wls/applications/sql_predict/apps              ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
stopsignal=QUIT               ; signal used to kill process (default TERM)
user=dbcm                   ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/wls/applications/sql_predict/logs/pallas-celery.log        ; stdout log path, NONE for none; default AUTO
#environment=APP_ENV="ProductionAdmin"       ; process environment additions (def no adds)
stopasgroup=true
killasgroup=true
